# scripts/dockerpod-mongo.service
[Unit]
Description=mongo docker pod
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/local/bin/docker-compose -f /home/charles/codes/docker/pod-mongo/docker-compose.yml up
ExecStop=/usr/local/bin/docker-compose  -f /home/charles/codes/docker/pod-mongo/docker-compose.yml stop

[Install]
WantedBy=default.target

